# conf/supervisord.conf
[unix_http_server]
file=/dev/shm/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
user=root		     ;

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock ; use a unix:// URL  for a unix socket

; Workers light & heavy

[program:worker-light]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work light --queue=light,default --sleep=1 --tries=3
autostart=true
autorestart=true
user=nginx
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker-light.log

[program:worker-heavy]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work heavy --queue=heavy --sleep=3 --tries=1
autostart=true
autorestart=true
user=nginx
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker-heavy.log

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24
autorestart=true
user=root
priority=100

[program:udev]
command=/sbin/udevd
user=root
priority=100

[program:xfce4-session]
command=dbus-launch /usr/bin/xfce4-session --display=:1
user=root
autorestart=true
priority=200

[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
user=root
autorestart=true
priority=200

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
user=root
autorestart=true
priority=200

[program:x11vnc]
command=/usr/bin/x11vnc -xkb -noxrecord -noxfixes -noxdamage -display :1 -nopw -wait 5 -shared
user=root
autorestart=true
priority=200
